# config/immortalwrt-openwrt-24.10/files/etc/uci-defaults/99-remove-minidlna-js
#!/bin/sh

# Remove the miniDLNA JS status file
rm -f /www/luci-static/resources/view/status/include/80_minidlna.js

exit 0
